# cron/fetchprice.cron
*/15 * * * * /usr/bin/php /var/www/litecoinblockhalf.com/getprice.php